# site_integrations/diamant/supervisor/diamant-gold-price-worker.conf
[program:diamant-gold-price-worker]
command=/usr/bin/php /var/www/u0861209/data/www/diamant.uz/system/cron/gold_price_worker.php
directory=/var/www/u0861209/data/www/diamant.uz
autostart=true
autorestart=true
startsecs=5
startretries=10
user=u0861209
redirect_stderr=true
stdout_logfile=/var/www/u0861209/data/logs/diamant-gold-price-worker.log
environment=RABBITMQ_HOST="YOUR_VPS_IP",RABBITMQ_PORT="5672",RABBITMQ_USER="jewelry_user",RABBITMQ_PASSWORD="my_secure_password",RABBITMQ_QUEUE="gold_price_events"
