# systemd/kurs-sync.service
[Unit]
Description=Kurs site sync worker
After=network-online.target mysql.service mariadb.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Documents/kurs_update
EnvironmentFile=/home/ubuntu/Documents/kurs_update/.env
ExecStart=/home/ubuntu/Documents/kurs_update/.venv/bin/python /home/ubuntu/Documents/kurs_update/sync_worker.py
Restart=always
RestartSec=5
User=ubuntu

[Install]
WantedBy=multi-user.target
